# musicpolice-web.service
[Unit]
Description=MusicPolice Web Interface
Documentation=https://github.com/your-repo/musicpolice
After=network.target musicpolice.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/musicpolice
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/musicpolice/venv/bin/python3 /home/pi/musicpolice/app.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=musicpolice-web

# Resource limits
MemoryMax=128M

[Install]
WantedBy=multi-user.target
